# deploy/launchd/com.octravpn.node.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.octravpn.node</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/octravpn-node</string>
        <string>--config</string>
        <string>/usr/local/etc/octravpn/node.toml</string>
        <string>run</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- macOS utun requires root. -->
    <key>UserName</key>
    <string>root</string>
    <key>GroupName</key>
    <string>wheel</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/octravpn-node.out.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/octravpn-node.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
</dict>
</plist>
